# deploy/api.service
[Unit]
Description=Crypto Trader API Server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/projects/crypto-trader
EnvironmentFile=/root/projects/crypto-trader/api/.env
Environment=PYTHONPATH=/root/projects/crypto-trader
Environment=PIPENV_PIPFILE=/root/projects/crypto-trader/Pipfile

ExecStart=pipenv run uvicorn \
    api.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 2 \
    --log-level info

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
