# deploy/dashboard-pro.service
[Unit]
Description=Crypto Trader Dashboard Pro (Next.js SSR)
After=network.target api.service
Wants=api.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/projects/crypto-trader/dashboard-pro
EnvironmentFile=/root/projects/crypto-trader/dashboard-pro/.env.local
Environment=PORT=3001

ExecStart=/usr/bin/pnpm start

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
